# cryptobot.service
[Unit]
Description=a proprietary cryptocurrency trading bot

[Service]
ExecStart=/usr/local/bin/cryptobot analyze-market

