# cryptobot.timer
[Unit]
Description=Analyze the network and make trades every hour

[Timer]
OnCalendar=*:30
Persistent=false

[Install]
timers.target
