# cryptobot-check-price.timer
[Unit]
Description=save cryptocurrency prices every minute

[Timer]
OnActiveSec=*:*
Persistent=false

[Install]
WantedBy=timers.target
