# cryptobot-check-price.service
[Unit]
Description=save cryptocurrency prices every minute

[Service]
ExecStart=/usr/local/bin/cryptobot check-price
User=cryptobot
